# src/conf/supervisor-ckan-server.conf
; =======================================================
; Supervisor configuration for CKAN background job worker
; =======================================================

; 1. Copy this file to /etc/supervisr/conf.d
; 2. Make sure the paths below match your setup


[program:ckan-server]

; Use the full paths to the virtualenv and your configuration file here.
command=/usr/lib/ckan/default/bin/paster serve --reload /etc/ckan/default/ckan.ini


; User the worker runs as.
user=root


; Start just a single worker. Increase this number if you have many or
; particularly long running background jobs.
; numprocs=2
; process_name=%(program_name)s-%(process_num)02d


; Log files.
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0


; Make sure that the worker is started on system start and automatically
; restarted if it crashes unexpectedly.
autostart=true
autorestart=true


; Number of seconds the process has to run before it is considered to have
; started successfully.
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 6000


[eventlistener:stdout]
command = supervisor_stdout
buffer_size = 100
events = PROCESS_LOG
result_handler = supervisor_stdout:event_handler
